# DEBIAN/postinst
#!/bin/sh
# postinst script for awslogs-agent

set -e

# Have awslogs-agent start on boot
update-rc.d awslogs-agent defaults 50
